# hermes-core/skills/software-development/hermes-android-gui/templates/hermes-gateway.service
[Unit]
Description=Hermes Native Gateway (full host access)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/user
Environment=HERMES_HOME=/home/user/.hermes-native-gateway
Environment=LITELLM_API_KEY=*** Environment=GLM_API_KEY=*** Environment=HERMES_API_SERVER_PORT=8643
ExecStart=/home/user/.hermes/hermes-agent/venv/bin/hermes gateway run
Restart=on-failure
RestartSec=5
KillMode=process
Delegate=yes

[Install]
WantedBy=default.target
